# scripts/lint_watcher_task.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  Aufgabenplanung-Definition fuer den Auto-Linter (Phase 2)
  =========================================================

  Importieren ueber PowerShell (als Administrator NICHT noetig):
      schtasks /Create /XML scripts\lint_watcher_task.xml /TN "EnglischLernzentrumLinter"

  Anschliessend startet der Linter automatisch beim naechsten Login.

  WICHTIG: Vor dem Import den absoluten Pfad zur start_lint_watcher.ps1 anpassen!
  Suche dazu im Element <Command>powershell.exe</Command> nach <Arguments>
  und ersetze C:\Users\thoma\... mit dem tatsaechlichen Pfad zu deiner
  englisch-lernzentrum/scripts/start_lint_watcher.ps1.
-->
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Startet den Auto-Linter fuer das Englisch-Lernzentrum bei Anmeldung.</Description>
    <Author>Thomas Porsche</Author>
    <URI>\EnglischLernzentrumLinter</URI>
  </RegistrationInfo>
  <Triggers>
    <LogonTrigger>
      <Enabled>true</Enabled>
    </LogonTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <DisallowStartOnRemoteAppSession>false</DisallowStartOnRemoteAppSession>
    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>powershell.exe</Command>
      <Arguments>-WindowStyle Hidden -ExecutionPolicy Bypass -File "K:\OneDrive - Hellweg-Realschule\01 - Schule\01 - Englisch\englisch-lernzentrum\scripts\start_lint_watcher.ps1"</Arguments>
    </Exec>
  </Actions>
</Task>
